# 12/rootfs/etc/crontabs/nextcloud-cron-15min
*/15 * * * * su-exec www-data php7 -f /var/www/cron.php
