# 15/rootfs/etc/crontabs/nextcloud-cron-15min
*/15 * * * * su-exec www-data php -f /var/www/cron.php
