# Dockerfiles/rootfs/etc/crontabs/nextcloud-cron-15min
*/15 * * * * su-exec root php -f /var/www/cron.php
