# Dockerfiles/rootfs/etc/crontabs/nextcloud-cron-15min
# Console has to be executed with the same user as the web server is operated
*/15 * * * * su-exec www-data php -f /var/www/cron.php
